# systemd/infomagic-touch.service
[Unit]
Description=InfoMagic Touch Display
After=network-online.target infomagic-backend.service

[Service]
Type=simple
User=@APP_USER@
SupplementaryGroups=video render
Environment=DISPLAY=:0
ExecStartPre=/bin/bash -c 'if ! pgrep -x Xorg > /dev/null; then /usr/bin/Xorg :0 -nolisten tcp vt7 & sleep 2; fi'
ExecStart=/usr/bin/chromium --kiosk --password-store=basic --use-mock-keychain --window-position=0,0 http://localhost:3000/ui/touch.html
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target