# systemd/infomagic-tv.service
[Unit]
Description=InfoMagic TV Display
After=network-online.target infomagic-backend.service infomagic-x.service
Requires=infomagic-x.service

[Service]
Type=simple
User=@APP_USER@
Environment=DISPLAY=:0
ExecStart=/usr/bin/chromium \
  --kiosk \
  --window-position=0,0 \
  --window-size=1920,1080 \
  --noerrdialogs \
  --disable-infobars \
  --disable-session-crashed-bubble \
  --disable-translate \
  --disable-features=TranslateUI \
  http://localhost:3000/ui/tv.html
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target