# systemd/weston.service
[Unit]
Description=Weston Wayland Compositor (Kiosk)
After=systemd-logind.service
Requires=systemd-logind.service

[Service]
Type=simple
User=root
Environment=XDG_RUNTIME_DIR=/run/weston
ExecStartPre=/bin/mkdir -p /run/weston
ExecStartPre=/bin/chmod 0700 /run/weston
ExecStart=/usr/bin/weston \
  --backend=drm-backend.so \
  --shell=kiosk-shell.so
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target

